# stage2/98-app/files/mpc-app.service
[Unit]
Description=mpc-app service
After=multi-user.target

[Service]
Type=idle
ExecStart=/home/mpc/run-app.sh

[Install]
WantedBy=multi-user.target
